# apps/runner/packaging/deb/DEBIAN/prerm
#!/bin/bash
set -e

# Stop the service if it's running
if systemctl is-active --quiet daytona-runner.service; then
    systemctl stop daytona-runner.service
fi

# Disable the service
if systemctl is-enabled --quiet daytona-runner.service; then
    systemctl disable daytona-runner.service
fi

exit 0
